# deploy/backup/tita-backup.service
[Unit]
Description=Tita PostgreSQL backup → S3
Wants=tita-backup.timer
After=network-online.target postgresql@15-main.service

[Service]
Type=oneshot
User=postgres
Group=postgres
EnvironmentFile=/etc/tita-api/backup.env
ExecStart=/usr/local/bin/tita-backup.sh
StandardOutput=append:/var/log/tita-api/backup.log
StandardError=append:/var/log/tita-api/backup.log
# Failure → Telegram alert qua OnFailure unit (task 09)
OnFailure=tita-alert@%n.service

# Resource limit — backup không được ăn hết RAM box.
MemoryMax=512M
CPUWeight=50
IOWeight=50
