# deploy/monitoring/tita-health.service
[Unit]
Description=Tita local health check

[Service]
Type=oneshot
EnvironmentFile=/etc/tita-api/alert.env
ExecStart=/usr/local/bin/tita-health-check.sh
